# build/mazezoom_supervisord.conf
[program:channel_realtime]
command=python /opt/mazezoom/src/mazezoom/creepers/channel_realtime_worker.py 
autorstart=true 
stdout_logfile=/data/log/channel_realtime.log

[program:channel_schedule]
command=python /opt/mazezoom/src/mazezoom/creepers/channel_schedule.py 
autorstart=true 
stdout_logfile=/data/log/channel_schedule.log

[program:dispatcher]
command=python /opt/mazezoom/src/mazezoom/creepers/position_dispatcher.py
stdout_logfile=/data/log/supervisord/position.log
stderr_logfile=/data/log/supervisord/position.log
autostart=true
autorestart=true

[program:schedule]
command=python /opt/mazezoom/src/mazezoom/creepers/position_schedule.py
stdout_logfile=/data/log/supervisord/pschedule.log
stderr_logfile=/data/log/supervisord/pschedule.log
autostart=true
autorestart=true
