# build/mazezoom_supervisord.conf
[program:dispatcher]
command=python /home/ubuntu/mazezoom/src/mazezoom/creepers/position_dispatcher.py 
autorstart=true 
stdout_logfile=/data/log/poisiton.log

[program:dispatcher_schedule]
command=python /home/ubuntu/mazezoom/src/mazezoom/creepers/position_schedue.py 
autorstart=true 
stdout_logfile=/data/log/schedule.log

[program:channel_realtime]
command=python /home/ubuntu/mazezoom/src/mazezoom/creepers/channel_realtime_worker.py 
autorstart=true 
stdout_logfile=/data/log/channel_realtime.log
